# packaging/postinst
#!/bin/sh

systemctl enable graphite-ch-optimizer.service
systemctl start graphite-ch-optimizer.service
