# packaging/prerm
#!/bin/sh

systemctl disable graphite-ch-optimizer.service
systemctl stop graphite-ch-optimizer.service
